# wings-cgroups/v1-legacy/t0-host-baseline/t0c-property-reconciler/wings-cg-reconcile.service
[Unit]
Description=Reconcile wings cgroup resource properties (slices + docker scopes)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wings-cg-reconcile.sh
