# pwmcp/containers/pwmcp/supervisord.shared.conf
; pwmcp shared-browser-mode supervisor configuration
; ====================================================
; Selected instead of pwmcp.conf when PWMCP_BROWSER_MODE=shared. Runs ONE
; persistent headless Chromium (program:chromium) that mcp / devtools-mcp
; attach to over CDP, instead of each launching its own browser. Manages
; six programs:
;   run-server      — Playwright native WebSocket server on :3000 (unchanged)
;   chromium        — the single shared Chromium, CDP on 127.0.0.1:9222 only
;   mcp             — @playwright/mcp attached via --cdp-endpoint, :8931
;   devtools-mcp    — chrome-devtools-mcp attached via --browser-url, :8932
;   lighthouse-mcp  — UNCHANGED: keeps its own per-audit browser launch
;                     (out of scope for shared-mode wiring per P03 handoff)
;   admin-server    — small HTTP admin API on 0.0.0.0:%(ENV_PWMCP_ADMIN_PORT)s,
;                     internal Docker network only (never published, never
;                     routed through Traefik — see ciu.compose.yml.j2).
;
; Safeguards implemented here:
;   1. autorestart=true + startretries on [program:chromium] — crash-restart.
;      Attached MCP servers reconnect on the NEXT tool call (they dial the
;      CDP endpoint per-call/session, not once at process start) rather than
;      caching a dead connection; wait-for-cdp.sh below is what bounds the
;      analogous wait during their own startup.
;   2. admin-server owns /browser/reset (close contexts) and /browser/restart
;      (supervisord-mediated hard restart) — session/state reset on demand.
;   3. State-bleed containment: NOT mitigated by design (self-review,
;      2026-07-13) -- --isolated was originally applied here but broke the
;      cross-tool workflow (DevTools traced a different browser context than
;      the one Playwright navigated). It has been removed; mcp and
;      devtools-mcp share Chromium's one default context on purpose so the
;      cross-tool workflow works. Cookies/storage bleed across concurrent
;      sessions on the shared browser as a result -- accepted, documented
;      residual risk (see docs/SECURITY.md "State-bleed residual"), not a
;      safeguard. Use per-session mode when per-consumer isolation matters.
;   4. admin-server also runs the optional idle-recycle loop
;      (PWMCP_BROWSER_MAX_IDLE_S).
;   5. Chromium's --remote-debugging-port binds to 127.0.0.1 only, inside
;      this container; it is never published, never on the compose ports
;      list, never given a Traefik route (see ciu.compose.yml.j2 / SECURITY.md).
;
; Startup ordering: wait-for-cdp.sh bounds how long mcp/devtools-mcp wait for
; the CDP port before exec'ing the real attach command (see that script for
; the "never wedge, let autorestart take over" rationale).

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:run-server]
command=playwright run-server --port 3000 --host 0.0.0.0
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"
priority=1

[program:chromium]
; The single persistent shared browser. CDP binds to 127.0.0.1 ONLY (never
; 0.0.0.0) — safeguard 5 (CDP never leaves the container); the Docker
; network cannot reach 127.0.0.1 of another container, so this is
; sufficient even without an extra firewall rule. --user-data-dir under
; /tmp keeps the profile in the world-writable runtime area used elsewhere.
command=%(ENV_PWMCP_CHROMIUM_PATH)s --headless=new --remote-debugging-port=%(ENV_PWMCP_CDP_PORT)s --remote-debugging-address=127.0.0.1 --no-sandbox --disable-setuid-sandbox --disable-gpu --user-data-dir=/tmp/pwmcp-shared-chromium-profile
autorestart=true
startretries=50
startsecs=1
; autorestart=true + startretries=50 gives a bounded-but-generous crash-restart
; window (supervisord's default backoff is capped); mechanism test in
; scripts/smoke-endpoints.sh shared-mode pass: kill -9 this program's pid and
; assert supervisord reports RUNNING again within a bounded poll window.
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"
priority=5

[program:mcp]
; Attaches to the shared browser instead of launching its own. Self-review
; (2026-07-13) found that passing --isolated here AND on devtools-mcp put
; each attached server in its own incognito-style browser context, which
; silently defeated the cross-tool workflow this whole package exists for
; (Playwright navigates a page in ITS context; DevTools traces a blank
; chrome://new-tab-page/ in a DIFFERENT context — verified via the smoke
; script's cross-tool CHECK, see P03-LOG.md and P03-SELFREVIEW.md). --isolated
; is removed here so both attached servers share Chromium's one default
; browser context, matching the handoff's core goal. State-bleed containment
; (safeguard 3) is downgraded from "mitigated via --isolated" to an accepted,
; explicitly documented residual risk (see docs/SECURITY.md) — sessions on
; the shared browser DO see each other's cookies/DOM state on the same
; origin/context; use per-session mode when that isolation is required.
; wait-for-cdp.sh bounds the wait for chromium's CDP port to come up
; (startup-ordering contract) before exec'ing playwright-mcp; if the deadline
; is hit anyway, playwright-mcp's own connection failure triggers supervisord
; autorestart, so the program never wedges half-attached.
command=/usr/local/bin/wait-for-cdp.sh 127.0.0.1 %(ENV_PWMCP_CDP_PORT)s 30 -- playwright-mcp --cdp-endpoint http://127.0.0.1:%(ENV_PWMCP_CDP_PORT)s --port 8931 --host 0.0.0.0 --allowed-hosts %(ENV_PWMCP_MCP_ALLOWED_HOSTS)s
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"
priority=10

[program:devtools-mcp]
; Attaches to the shared browser via --browser-url instead of launching its
; own (chrome-devtools-mcp's own-launch flags no longer apply — the browser
; already exists). No --isolated (see [program:mcp] above): shares the
; default context with mcp so the cross-tool workflow (drive a page with
; Playwright, profile it with DevTools) actually lands on the same page.
command=/usr/local/bin/wait-for-cdp.sh 127.0.0.1 %(ENV_PWMCP_CDP_PORT)s 30 -- mcp-proxy --port 8932 -- chrome-devtools-mcp --browser-url http://127.0.0.1:%(ENV_PWMCP_CDP_PORT)s
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"
priority=10

[program:lighthouse-mcp]
; UNCHANGED from per-session mode: Lighthouse keeps its own per-audit
; Chromium launch (chrome-launcher), NOT wired to the shared browser. This
; is explicit P03 scope (see handoff "Out Of Scope").
command=mcp-proxy --port 8933 -- node /opt/pwmcp/lighthouse-mcp/index.js
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"
priority=10

[program:admin-server]
; Small stdlib-only Node HTTP admin API. Binds 0.0.0.0:%(ENV_PWMCP_ADMIN_PORT)s
; (the compose template never publishes this port and never adds a Traefik
; route for it — internal Docker network only, safeguard 2). Talks to
; supervisord over the same unix socket as supervisorctl for
; program-restart, and to Chromium's CDP HTTP/WS endpoint for
; reset/health/idle-recycle. Grants itself no capability beyond
; context-close / program-restart / health, per the handoff.
command=node /opt/pwmcp/admin-server/index.js
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu",PWMCP_ADMIN_PORT="%(ENV_PWMCP_ADMIN_PORT)s",PWMCP_CDP_PORT="%(ENV_PWMCP_CDP_PORT)s",PWMCP_BROWSER_MAX_IDLE_S="%(ENV_PWMCP_BROWSER_MAX_IDLE_S)s",PWMCP_SUPERVISOR_SOCK="/tmp/supervisor.sock"
priority=10
