# pwmcp/containers/pwmcp/supervisord.conf
; pwmcp unified supervisor configuration
; =======================================
; Manages three services inside the pwmcp container:
;   run-server    — Playwright native WebSocket server on :3000
;   mcp           — @playwright/mcp HTTP/SSE server on :8931
;   devtools-mcp  — chrome-devtools-mcp via mcp-proxy on :8932
;
; Designed to run as uid 1000 (non-root). All paths use /tmp or /var/log/supervisor
; so no root-owned directories are required at runtime.
;
; supervisord is launched as PID 1 with --nodaemon (see Dockerfile ENTRYPOINT).
; It reaps zombie children and forwards SIGTERM to all programs on shutdown.
;
; PWMCP_MCP_ALLOWED_HOSTS env var controls the @playwright/mcp --allowed-hosts
; value. Defaults to "localhost:8931,127.0.0.1:8931" for standalone usage.
; ciu compose injects the deployment-specific container names at runtime, e.g.:
;   PWMCP_MCP_ALLOWED_HOSTS=pwmcp:8931,myproject-dev-pwmcp:8931
;
; PWMCP_DEVTOOLS_ALLOWED_HOSTS tracks the expected hosts for chrome-devtools-mcp.
; mcp-proxy does not have a native --allowed-hosts flag, so the env var is for
; documentation and external-mode Traefik rules.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

; No unix_http_server block — we do not need supervisorctl from inside the
; container. Omitting it avoids the default /var/run/supervisor.sock path
; that requires root.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:run-server]
command=playwright run-server --port 3000 --host 0.0.0.0
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
; Propagate stdout/stderr to the parent process (docker logs).
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"

[program:mcp]
; playwright-mcp is the bin installed by `npm install -g @playwright/mcp`.
; It symlinks to /usr/lib/node_modules/@playwright/mcp/cli.js.
; --allowed-hosts is injected at container start via PWMCP_MCP_ALLOWED_HOSTS.
; --executable-path bypasses @playwright/mcp's bundled playwright-core chromium
; discovery (which may target a different revision) and uses the binary from the
; playwright JS package already installed in the image. PWMCP_CHROMIUM_PATH is
; exported by the entrypoint script from /etc/pwmcp-chromium-path.txt (baked at
; image build time).
command=playwright-mcp --headless --browser chromium --no-sandbox --executable-path %(ENV_PWMCP_CHROMIUM_PATH)s --port 8931 --host 0.0.0.0 --allowed-hosts %(ENV_PWMCP_MCP_ALLOWED_HOSTS)s
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"

[program:devtools-mcp]
; chrome-devtools-mcp is stdio-only (no --port flag); mcp-proxy wraps it and
; exposes the MCP protocol over HTTP/SSE on port 8932.
; --executable-path uses the same baked chromium path as @playwright/mcp.
; --isolated creates a temporary user-data-dir cleaned up on exit.
; --no-sandbox matches the hardening pattern from the existing [program:mcp].
; PWMCP_CHROMIUM_PATH is exported by the entrypoint script.
command=mcp-proxy --command "chrome-devtools-mcp --headless --no-sandbox --executable-path %(ENV_PWMCP_CHROMIUM_PATH)s --isolated" --port 8932
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=HOME="/home/ubuntu"
