# scripts/gstammtisch-guide/files/etc/systemd/system/dev-workloads.slice
[Unit]
Description=Development workloads (dev containers, builds) — low priority, OOM-first
Before=slices.target

[Slice]
# --- Memory ---
# Throttle reclaim (soft) at 8G; hard cap at 14G. OOM fires INSIDE this slice,
# never touching the game server in its own cgroup.
MemoryHigh=8G
MemoryMax=14G

# --- CPU / IO weight ---
# Half the default weight; Soulmask (default 100) gets 2:1 under contention.
CPUWeight=50
IOWeight=50

# --- Build I/O cap ---
# Stop docker build/pull from saturating the disk. Adjust device if not /dev/vda.
IOReadBandwidthMax=/dev/vda 200M
IOWriteBandwidthMax=/dev/vda 200M

# --- systemd-oomd: kill the worst offender in THIS slice first under pressure ---
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=60%
ManagedOOMSwap=kill

# NOTE: memory.zswap.writeback=1 for this slice (allow eviction of dev pages to
# disk) is applied by setup-cgroups.sh — systemd has no property for it.
# Launch dev containers into this slice:
#   docker run --cgroup-parent=dev-workloads.slice --label workload=dev ...
# or in compose:  cgroup_parent: dev-workloads.slice
