# scripts/gstammtisch-guide/files/etc/systemd/system/soulmask-cgroup-watcher.service
[Unit]
Description=Reconcile loop: apply cgroup knobs to every ready Soulmask instance (game-log readiness, not RCON)
Documentation=file:///home/vb/volkb79-2/vbpub/scripts/gstammtisch-guide/SOULMASK.md
After=docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/soulmask-cgroup-watcher.sh
Restart=always
RestartSec=10
# Allow slow startup (RCON can take several minutes on first cold-disk load)
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
