# scripts/gstammtisch-guide/files/etc/systemd/system/soulmask-cgroup-watcher.service
[Unit]
Description=Wait for Soulmask RCON ready then apply cgroup knobs
After=docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/soulmask-cgroup-watcher.sh
Restart=always
RestartSec=10
# Allow slow startup (RCON can take several minutes on first cold-disk load)
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
