# scripts/gstammtisch-guide/files/etc/systemd/system/soulmask-pak-ramdisk.service
[Unit]
Description=Soulmask pak file ramdisk (shared tmpfs, bind-mounted into every PAK_RAMDISK=1 instance)
Documentation=file:///home/vb/volkb79-2/vbpub/scripts/gstammtisch-guide/SOULMASK.md
# Must run before Docker so containers mount the tmpfs-backed Paks directory.
Before=docker.service
After=local-fs.target
DefaultDependencies=no

[Service]
# Run inside soulmask-paks.slice so pak tmpfs pages are charged to that slice,
# not the root cgroup. Gives independent memory.min + writeback control for pak.
Slice=soulmask-paks.slice
Type=oneshot
RemainAfterExit=yes
# N-instance (2026-07-07): enumerates /etc/gstammtisch/instances.d/*.env for
# PAK_RAMDISK=1 instances and bind-mounts the ONE shared tmpfs into each of
# their volumes — see SOULMASK.md "Multi-instance operations".
ExecStart=/usr/local/sbin/soulmask-pak-ramdisk-setup.sh
# Teardown: unmount every bind target recorded in the setup script's state
# file, then the shared tmpfs itself. A dedicated script now (not inline
# /bin/sh -c) because the state file lists a variable number of targets.
ExecStop=/usr/local/sbin/soulmask-pak-ramdisk-teardown.sh
# Allow up to 2 min for large-pak copy on slow cold disk.
# (Pre-existing bug fixed in passing: a trailing inline comment on the same
# line as TimeoutStartSec= is not valid systemd unit syntax and was silently
# ignored — `systemd-analyze verify` flags it.)
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
