# scripts/gstammtisch-guide/files/etc/systemd/system/soulmask-pak-ramdisk.service
[Unit]
Description=Soulmask pak file ramdisk (tmpfs over Paks directory)
Documentation=file:///home/vb/volkb79-2/vbpub/scripts/gstammtisch-guide/SOULMASK.md
# Must run before Docker so the container mounts the tmpfs-backed Paks directory.
Before=docker.service
After=local-fs.target
DefaultDependencies=no

[Service]
# Run inside soulmask-paks.slice so pak tmpfs pages are charged to that slice,
# not the root cgroup. Gives independent memory.min + writeback control for pak.
Slice=soulmask-paks.slice
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/soulmask-pak-ramdisk-setup.sh
# Teardown: read PAK_DIR from state file written by the setup script.
# Unmount in reverse order: bind first, then the backing tmpfs.
ExecStop=/bin/sh -c '\
  STATE=/run/soulmask-pak-ramdisk.state; \
  PAK_DIR=""; \
  [ -f "$STATE" ] && . "$STATE"; \
  [ -n "$PAK_DIR" ] && umount "$PAK_DIR" 2>/dev/null || true; \
  umount /mnt/soulmask-paks 2>/dev/null || true; \
  rm -f "$STATE"'
TimeoutStartSec=120   # allow up to 2 min for large-pak copy on slow cold disk
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
