# modern-debian-tools-python-debug/host-setup/units/mdt-host-slices.service
# mdt host-setup — applies the runtime half of the dev-tier governance:
# besteffort.slice measured IO caps, interactive zswap policy, per-container
# caps for bench/buildkit/devcontainer containers. Idempotent oneshot; run at
# boot (enabled) and periodically via mdt-host-slices.timer.
[Unit]
Description=mdt host-setup — runtime dev-tier IO caps + per-container limits
Documentation=file:/etc/mdt/host-setup.env
# The container sweep needs docker; the slice property calls do not — the
# script degrades gracefully when docker is down.
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mdt-apply-dev-caps.sh

[Install]
WantedBy=multi-user.target
