# groop/src/groop/assets/systemd/groop.service
[Unit]
Description=groop read-only daemon
After=network.target

[Service]
Type=simple
User=root
Group=groop
UMask=0007
RuntimeDirectory=groop
RuntimeDirectoryMode=0750
ExecStart=/usr/bin/env groop daemon serve --socket /run/groop/groop.sock
Restart=on-failure

# Create the groop group and add approved users before enabling this unit.
# The daemon remains read-only; clients attach over the group-readable socket.

[Install]
WantedBy=multi-user.target
